# rootfs/etc/s6/alerts/run
#!/bin/sh

cd /app

exec python -m src.main 2>&1
