# rootfs/etc/s6/cron/run
#!/bin/sh

cd /app

exec /usr/sbin/cron -f 2>&1
